# .openshift/cron/daily/notify.cron
#!/usr/bin/env bash
curl -G http://tenders.rudra.com.np/cron/notification


